# howdy-authselect/howdy-authselect.path
[Unit]
Description=Watch for authselect configuration changes
Documentation=man:howdy(1)

[Path]
PathChanged=/etc/authselect/system-auth
PathChanged=/etc/authselect/password-auth

[Install]
WantedBy=multi-user.target
